# scripts/launchd/com.user.ai-weekly-news.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.ai-weekly-news</string>

  <key>WorkingDirectory</key>
  <string>/Users/angushyx/Desktop/ai-daily-news</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/angushyx/.nvm/versions/node/v18.18.2/bin/node</string>
    <string>/Users/angushyx/Desktop/ai-daily-news/src/index.js</string>
    <string>--mode=weekly</string>
  </array>

  <!-- 每週一 10:00 (本機時區)；Weekday: 0=Sun, 1=Mon ... 6=Sat -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>1</integer>
    <key>Hour</key><integer>10</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/angushyx/.local/bin:/Users/angushyx/.nvm/versions/node/v18.18.2/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/ai-weekly-news.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/ai-weekly-news.err</string>
</dict>
</plist>
